# scripts/ansible/roles/install_nodeos_services/files/keosd.service
[Unit]
Description= Keosd eos wallet
After=network.target

[Service]
Type=simple
ExecStart=keosd
TimeoutStartSec=0

[Install]
WantedBy=default.target